# contrib/systemd/mnemos-sync.service
[Unit]
Description=Mnemos — auto-cron federation bridge (#104) sync-peers.sh
Documentation=file:///usr/share/doc/mnemos/admin/ssh-sync-hardening.md
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=mnemos-sync
Group=mnemos-sync
EnvironmentFile=/etc/mnemos/sync.env
ExecStart=/usr/local/sbin/sync-peers.sh
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
NoNewPrivileges=yes
ReadWritePaths=/var/lib/mnemos-sync /var/log/mnemos-sync.log
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target