# contrib/systemd/mnemos-watcher.service
[Unit]
Description=Mnemos — File Watcher
After=network.target

[Service]
Type=simple
ExecStart=%h/LABs/AI/mnemos/.venv/bin/mnemos watch
Restart=on-failure
RestartSec=10
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=default.target
