# contrib/systemd/mnemos-server.service
[Unit]
Description=Mnemos — Memory & Knowledge Server (HTTP API)
After=network.target

[Service]
Type=simple
ExecStart=%h/LABs/AI/mnemos/.venv/bin/mnemos serve
Restart=on-failure
RestartSec=5
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=default.target
